# src/supervisor/supervisord.conf
[supervisord]
pidfile = /var/run/supervisord.pid
nodaemon=true # arkaplanda surekli olarak calisacak

[program:cron] 
command=cron -f
autostart=true
autorestart=true